# scripts/systemd/signal-works-morning.service
[Unit]
Description=agentsHQ Signal Works + Catalyst Works morning runner (runs inside orc-crewai container)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker exec orc-crewai python3 -m signal_works.morning_runner
StandardOutput=append:/var/log/signal_works_morning.log
StandardError=append:/var/log/signal_works_morning.log
TimeoutStartSec=240min

[Install]
WantedBy=multi-user.target
