# agents/outfit_stylist/outfit-stylist.service
[Unit]
Description=AminoAI — Aesthetic Intelligence
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/agentsHQ/agents/outfit_stylist
EnvironmentFile=/root/agentsHQ/.env
ExecStart=/usr/bin/python3 /root/agentsHQ/agents/outfit_stylist/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
